# data/units/clouseaud.socket
[Unit]
Description=Clouseau Service Sockets

[Socket]
ListenStream=22522

[Install]
WantedBy=sockets.target
